# curator/crontab
# Snapshot Every 1 hour
1 * * * * curator /actions/snapshot.yml

