# curator/crontab
# Snapshot Every midnight
* 0 * * * curator /actions/snapshot.yml

