# curator/crontab
# Snapshot Every 1 hour
0 0 * * * curator /actions/snapshot.yml

